# Configurazioni/Edificio A (Studi di Programmazione)/Router A/rc.local
#DIRECTORY = /etc/rc.local

### Configurazione interfacce

//uffici studio programmazione (locale)
ifconfig eth1 192.168.4.255 netmask 255.255.255.0 broadcast 192.168.4.255

//verso router (b)
ifconfig eth2 192.168.0.121 netmask 255.255.255.0 broadcast 192.168.0.255

